# roles/ctf_challenge_socket_backup/files/units/dk-tasktracker.socket
[Unit]
Description=DK tasktracker Socket
PartOf=dk-tasktracker@.service

[Socket]
ListenStream=0.0.0.0:9997
Accept=yes
MaxConnections=64
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
