# roles/ctf_challenge_socket_backup/files/units/en-FSBS.socket
[Unit]
Description=EN FSBS Socket
PartOf=en-FSBS@.service

[Socket]
ListenStream=0.0.0.0:9994
Accept=yes
MaxConnections=64
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
