# roles/ctf_challenge_socket_backup/files/units/nqh-insecure.socket
[Unit]
Description=NQH insecure Socket
PartOf=nqh-insecure@.service

[Socket]
ListenStream=0.0.0.0:9999
Accept=yes
MaxConnections=64
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
