# roles/ctf_challenge_socket/files/units/ctfsocket.slice
[Unit]
Description=Memory, CPU and IO control for ctf socket-activated daemons
Before=slices.target

[Slice]
MemoryAccounting=true
IOAccounting=true
CPUAccounting=true

CPUWeight=80
IOWeight=40
CPUQuota=80%

MemoryMax=40%
MemoryLimit=35%
BlockIOWeight=10
